# apps/api/crontab
0 * * * * cd /app && bun src/scripts/stripe-sync-entitlements.ts
